# .config/systemd/user/low-battery.timer
[Unit]
Description=Run low-battery every minute

[Timer]
OnCalendar=minutely

[Install]
WantedBy=wm.target
