# linux/.config/systemd/user/public-wifi.timer
[Unit]
Description=Run public-wifi every minute

[Timer]
OnCalendar=minutely

[Install]
WantedBy=wm.target
