# linux/.config/systemd/user/time-to-move.timer
[Unit]
Description=Run time to move every minute

[Timer]
OnCalendar=minutely

[Install]
WantedBy=wm.target
