# linux/systemd/system/tzupdate.timer
[Unit]
Description=Run tzupdate every minute

[Timer]
OnCalendar=minutely

[Install]
WantedBy=wm.target
