# linux/systemd/system/tzupdate.service
[Unit]
Description=Run tzupdate

[Service]
Type=oneshot
ExecStart=/bin/bash -c "%h/.local/bin/tzupdate"

[Install]
WantedBy=wm.target
