# linux/systemd/system/tzupdate.service
[Unit]
Description=Run tzupdate

[Service]
Type=oneshot
ExecStart=/usr/bin/tzupdate

[Install]
WantedBy=wm.target
